# profile.d/30-nvm.sh
if [ -e ~/.nvm/nvm.sh ]; then
  NVM_DIR=$(cd ~/.nvm/;pwd -P)
  source ~/.nvm/nvm.sh
fi
