# profile.d/brew
export HOMEBREW_GITHUB_API_TOKEN=d1dcbfb30f056865a76f2ce639cc0d28609f87ce