# profile.d/db_env.sh
/Users/vitaly/my/config/dotfiles/zsh/profile.d/db_env.sh