# profile.d/docker.sh
/Users/vitaly/my/config/dotfiles/zsh/profile.d/docker.sh